# config/memcached.conf
[program:memcached]
command=/usr/bin/memcached -p 11211 -u deploy -m 256 -c 1024 -t 4
